# etc/systemd/system/pendrive-attach@.service
[Unit]
Description=Automatically mount media devices

[Service]
Type=oneshot
ExecStart=/etc/automount/mountdrive.sh %i
